# mx-xfce/autostart/xfce4-clipman-plugin-autostart.desktop
[Desktop Entry]
Type=Application
Name=Clipman
Name[ar]=Clipman
Name[ast]=Clipman
Name[be]=Clipman
Name[bg]=Clipman
Name[ca]=Clipman
Name[cs]=Clipman
Name[da]=Clipman
Name[de]=Clipman
Name[el]=Clipman
Name[en_AU]=Clipman
Name[en_GB]=Clipman
Name[es]=Clipman
Name[eu]=Clipman
Name[fi]=Clipman
Name[fr]=Clipman
Name[gl]=Clipman
Name[he]=Clipman
Name[hr]=Clipman
Name[hu]=Clipman
Name[hy_AM]=Clipman
Name[id]=Clipman
Name[ie]=Clipman
Name[is]=Clipman
Name[it]=Clipman
Name[ja]=Clipman
Name[kk]=Clipman
Name[ko]=Clipman
Name[lt]=Clipman
Name[lv]=Clipman
Name[ms]=Clipman
Name[nb]=Clipman
Name[nl]=Clipman
Name[oc]=Clipman
Name[pl]=Clipman
Name[pt]=Clipman
Name[pt_BR]=Clipman
Name[ru]=Clipman
Name[sk]=Clipman
Name[sl]=Clipman
Name[sq]=Clipman
Name[sr]=Clipman
Name[sv]=Clipman
Name[th]=Clipman
Name[tr]=Clipman
Name[ug]=Clipman
Name[uk]=Clipman
Name[ur]=Clipman
Name[ur_PK]=Clipman
Name[vi]=Clipman
Name[zh_CN]=Clipman
Name[zh_TW]=Clipman
Comment=Clipboard manager
Comment[ar]=Clipboard manager
Comment[ast]=Clipboard manager
Comment[be]=Clipboard manager
Comment[bg]=Clipboard manager
Comment[ca]=Clipboard manager
Comment[cs]=Clipboard manager
Comment[da]=Clipboard manager
Comment[de]=Clipboard manager
Comment[el]=Clipboard manager
Comment[en_AU]=Clipboard manager
Comment[en_GB]=Clipboard manager
Comment[es]=Clipboard manager
Comment[eu]=Clipboard manager
Comment[fi]=Clipboard manager
Comment[fr]=Clipboard manager
Comment[gl]=Clipboard manager
Comment[he]=Clipboard manager
Comment[hr]=Clipboard manager
Comment[hu]=Clipboard manager
Comment[hy_AM]=Clipboard manager
Comment[id]=Clipboard manager
Comment[ie]=Clipboard manager
Comment[is]=Clipboard manager
Comment[it]=Clipboard manager
Comment[ja]=Clipboard manager
Comment[kk]=Clipboard manager
Comment[ko]=Clipboard manager
Comment[lt]=Clipboard manager
Comment[lv]=Clipboard manager
Comment[ms]=Clipboard manager
Comment[nb]=Clipboard manager
Comment[nl]=Clipboard manager
Comment[oc]=Clipboard manager
Comment[pl]=Clipboard manager
Comment[pt]=Clipboard manager
Comment[pt_BR]=Clipboard manager
Comment[ru]=Clipboard manager
Comment[sk]=Clipboard manager
Comment[sl]=Clipboard manager
Comment[sq]=Clipboard manager
Comment[sr]=Clipboard manager
Comment[sv]=Clipboard manager
Comment[th]=Clipboard manager
Comment[tr]=Clipboard manager
Comment[ug]=Clipboard manager
Comment[uk]=Clipboard manager
Comment[ur]=Clipboard manager
Comment[ur_PK]=Clipboard manager
Comment[vi]=Clipboard manager
Comment[zh_CN]=Clipboard manager
Comment[zh_TW]=Clipboard manager
TryExec=xfce4-clipman
Exec=xfce4-clipman
Icon=xfce4-clipman-plugin
OnlyShowIn=XFCE;X-Cinnamon;
Terminal=false
Hidden=true

